# deploy/systemd/trading-health.timer
[Unit]
Description=Run Trading Bot Daily Health Check

[Timer]
OnCalendar=*-*-* 08:10:00 America/New_York
Persistent=true
Unit=trading-health.service

[Install]
WantedBy=timers.target
