# deploy/systemd/trading-evolve-shadow.service
[Unit]
Description=Trading Bot Daily Shadow Mark + Live Drawdown Envelope Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=trader
WorkingDirectory=/home/trader/trading-bot
EnvironmentFile=/home/trader/trading-bot/.env
ExecStart=/home/trader/trading-bot/.venv/bin/python /home/trader/trading-bot/main.py evolve shadow-step
